# infra/systemd/portfolio-backup.timer
[Unit]
Description=Daily selickiy.space backup

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=20m
Unit=portfolio-backup.service

[Install]
WantedBy=timers.target
