# infra/systemd/portfolio.service
[Unit]
Description=Portfolio Next.js App
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=portfolio
Group=portfolio
WorkingDirectory=/var/www/portfolio
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/var/www/portfolio/.env
EnvironmentFile=-/var/lib/portfolio/release.env
StateDirectory=portfolio
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/www/portfolio/.next/cache /var/lib/portfolio
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
